# supervisor-configs/beachserver.ini
[program:beachserver]
user=ec2-user
autostart=true
autorestart=true
command=/home/ec2-user/codingbooth_env/codingbooth/env_run.sh python beachserver.py
