# supervisor-configs/codingbooth.ini
[program:codingbooth]
command = /usr/bin/uwsgi --socket 127.0.0.1:3031 -M -H /home/ec2-user/codingbooth_env -t 60 --need-app --module codingbooth:app --processes 2
user=ec2-user
stopsignal=INT
